# home/systemd/tier2-connect.service
[Unit]
Description=Tier-2 SSH tunnel to VPS (tun0 10.177.2.1/30)
Documentation=https://github.com/Cyrillicspb/vpn-infra/blob/main/docs/ARCHITECTURE.md
After=network-online.target
Wants=network-online.target
Before=vpn-routes.service watchdog.service
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=simple
User=root
EnvironmentFile=-/opt/vpn/.env
EnvironmentFile=-/run/vpn-active-backend.env
ExecStart=/opt/vpn/scripts/tier2-connect.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tier2-connect

[Install]
WantedBy=multi-user.target
