# home/watchdog/watchdog.service
[Unit]
Description=VPN Watchdog — центральный агент управления
After=network-online.target docker.service dnsmasq.service
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=root
WorkingDirectory=/opt/vpn/watchdog
ExecStart=/opt/vpn/watchdog/venv/bin/python watchdog.py
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
WatchdogSec=30
KillMode=process

# Логирование (watchdog пишет сам в LOG_FILE, stderr — для критических ошибок до старта)
StandardOutput=null
StandardError=append:/var/log/vpn-watchdog.log

# Окружение
EnvironmentFile=/opt/vpn/.env

[Install]
WantedBy=multi-user.target
